# fly/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:grpc-server]
command=python -m petstore_grpc
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:envoy]
command=/bin/sh -c "level=${ENVOY_LOG_LEVEL:-${LOG_LEVEL:-info}}; level=$(echo \"$level\" | tr '[:upper:]' '[:lower:]'); exec envoy -c /etc/envoy/envoy.yaml --log-level \"$level\""
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
